# Pyramid/dist/ro/var/spool/cron/crontabs/root
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.ECVuaS/crontab installed on Fri Jun  1 23:43:49 2007)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
#
# crontab for root
SHELL=/bin/bash
PATH=/usr/bin:/bin:/sbin:/usr/sbin:/usr/local/bin:/usr/local/sbin
#
# M H D M DOW	COMMAND
#
# Update WiFiAdmin rrds every 5 minutes
#
* * * * *	php /var/www/admin/create-update-rrds.php > /dev/null 2>&1
